# dist/systemd/etc/systemd/system/cacti-spine.service
[Unit]
Description=Cacti Spine Poller
After=network.target auditd.service

[Service]
ExecStart=/usr/local/spine/spine
KillMode=process
Restart=on-failure
